# deploy/petshop-live.service
[Unit]
Description=petshop-live demo storefront (tailnet-only)
After=network-online.target tailscaled.service

[Service]
WorkingDirectory=/var/tmp/nott/slack-work/petshop-live
Environment=HOST=100.65.70.81
Environment=PORT=4174
ExecStart=/usr/bin/node /var/tmp/nott/slack-work/petshop-live/server.js
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
